# scripts/sisters-whatsapp.service
[Unit]
Description=Sisters-On-WhatsApp Bot (WhatsApp Business API)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/Sisters-On-WhatsApp
Environment="PATH=/root/Sisters-On-WhatsApp/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=/root/Sisters-On-WhatsApp/.env
ExecStart=/root/Sisters-On-WhatsApp/venv/bin/python -m uvicorn src.whatsapp_webhook.server:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=10
StandardOutput=append:/var/log/sisters-whatsapp/access.log
StandardError=append:/var/log/sisters-whatsapp/error.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
